# root/etc/services.d/rso/run
#!/usr/bin/with-contenv bash

if [ -n "$SKIP_APP" ]; then
	exit 0
fi

nginx -g 'daemon off;'
